# cloudflared/rootfs/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the cloudflared service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
declare tunnel_token

# Get the 'tunnel_token' key from the user config options
tunnel_token=$(bashio::config 'tunnel_token')

# Check if tunnel token is provided
if [[ -z "${tunnel_token}" ]]; then
  bashio::log.fatal "No tunnel token configured!"
  bashio::log.fatal "Please configure a tunnel token in the add-on configuration."
  exit 1
fi

bashio::log.info "Starting Cloudflare Tunnel..."

# Create token file with secure permissions from the start (avoids race condition)
# Use /var/run which is typically more secure than /tmp
TOKEN_FILE="/var/run/cloudflared-token"
(umask 077; echo -n "${tunnel_token}" > "${TOKEN_FILE}")

# Clean up token file on exit or termination signals
trap 'rm -f "${TOKEN_FILE}"' EXIT INT TERM

# Run cloudflared with token file
exec cloudflared tunnel --no-autoupdate run --token-file "${TOKEN_FILE}"
